# data/indicator-sync.conf.in
description "Indicator Sync Service"

start on indicator-services-start
stop on desktop-end or indicator-services-end

respawn
respawn limit 2 10

pre-start script
	# The sync indicator shouldn't be on the greeter
	if [ "${USER}" = "lightdm" ]; then
		stop
	fi
end script

exec @pkglibexecdir@/indicator-sync-service
