# data/indicator-sync.conf.in
description "Indicator Sync Service"

# NOTE: Limiting only to Unity 7 right now as it's still using
# dbusmenu.  That can be lifted after it is ported to GMenu

start on indicator-services-start and xsession SESSION!=ubuntu-touch
stop on desktop-end or indicator-services-end

respawn
respawn limit 2 10

pre-start script
	# The sync indicator shouldn't be on the greeter
	if [ "${USER}" = "lightdm" ]; then
		stop
	fi
end script

exec @pkglibexecdir@/indicator-sync-service
